# deploy/systemd/medusa-maintenance.path
[Unit]
Description=Watch for user-approved Medusa maintenance requests

[Path]
PathExists=/opt/medusa/data/deploy/maintenance-request.json
Unit=medusa-maintenance.service

[Install]
WantedBy=multi-user.target
