# deploy/systemd/medusa-release-check.path
[Unit]
Description=Watch for Medusa release status check requests

[Path]
PathExists=/opt/medusa/data/deploy/release-check-request.json
Unit=medusa-release-check.service

[Install]
WantedBy=multi-user.target
