# deploy/systemd/medusa-release-apply.service
[Unit]
Description=Apply requested Medusa release upgrade
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/medusa
ExecStart=/opt/medusa/scripts/medusa-release-agent.py apply --repo /opt/medusa --data-dir /opt/medusa/data --compose-file docker-compose.yml --compose-file docker-compose.server.yml
