# deploy/systemd/medusa.service
[Unit]
Description=Medusa Docker Compose stack
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/medusa
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
ExecStart=/opt/medusa/scripts/medusa-release-agent.py compose --repo /opt/medusa --compose-file docker-compose.yml --compose-file docker-compose.server.yml -- up -d --build
ExecReload=/opt/medusa/scripts/medusa-release-agent.py compose --repo /opt/medusa --compose-file docker-compose.yml --compose-file docker-compose.server.yml -- up -d --build
ExecStop=/opt/medusa/scripts/medusa-release-agent.py compose --repo /opt/medusa --compose-file docker-compose.yml --compose-file docker-compose.server.yml -- stop
RemainAfterExit=yes
TimeoutStartSec=900
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
